# license-bot.service
[Unit]
Description=License Bot - Telegram Bot for license.gov.uz
After=network.target

[Service]
Type=simple
User=%I
WorkingDirectory=/home/%I/license_bot
Environment=PATH=/home/%I/license_bot/venv/bin
ExecStart=/home/%I/license_bot/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=license-bot

[Install]
WantedBy=multi-user.target
